# etc/rc.local
#!/bin/bash

ifconfig br0 192.168.12.1

/root/scripts/iptables.sh
sleep 3
/etc/init.d/monitorix restart


sleep 90
ifconfig | grep -v -e 'X ' | mail -s 'kjkao boot ifconfig' jeremy.kao@gmail.com

